# spotiplay.service
# Example systemd service file for Spotiplay Flask app
# Save as /etc/systemd/system/spotiplay.service and edit paths as needed

[Unit]
Description=Spotiplay Flask App (Gunicorn)
After=network.target

[Service]
WorkingDirectory=/home/ubuntu/spotiplay

ExecStart=/home/ubuntu/spotiplay/.venv/bin/gunicorn --workers 3 --bind 127.0.0.1:5000 app:app

Restart=always
RestartSec=5
TimeoutStopSec=30
SyslogIdentifier=spotiplay

[Install]
WantedBy=multi-user.target

